# tomcat.conf
[supervisord]
nodaemon=true

[program:tomcat]
command=/run.sh
